# crontab
*/1 * * * * /bin/bash -l -c 'cd /app && RAILS_ENV=production bundle exec rake newspaper:new_images --silent'
